# crontab/tankrain
DISPLAY=:0
PIG=hat
@reboot pi /usr/bin/python3.5 -m karmapi.tankrain /home/pi/karmapi/camera
